# startExpressVpn_Monitor.service
[Unit]
Description=ExpressVpn Monitoring

[Service]
Type=simple
ExecStart=/opt/Vpn_Monitor/ExpressVpn_Monitor.py 

[Install]
WantedBy=default.target
